# xmrig.timer
[Unit]
Description=Boot delay for XMRig Miner

[Timer]
OnBootSec=1min

[Install]
WantedBy=timers.target